# launchd/com.talktomepy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.talktomepy</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_DIR__/scripts/run_service.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_DIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/talktomepy.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/talktomepy.stderr.log</string>
</dict>
</plist>
